# meta-tellsticklove/recipes-support/ftdi-modprobe/files/99-usb-ftdi.rules
ACTION=="add", ATTRS{idVendor}=="1781", ATTRS{idProduct}=="0c31", RUN+="/sbin/modprobe ftdi_sio", RUN+="/sbin/modprobe usbserial"